# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10

[program:ttyd]
command=/usr/local/bin/ttyd -p 7681 -b /terminal -W bash -l
autostart=true
autorestart=true
priority=20

[program:dropbear]
command=/usr/sbin/dropbear -F -E -p 22 -R -B
autostart=true
autorestart=true
priority=30
startsecs=2
startretries=10

[program:bore]
command=/usr/local/bin/bore-wrapper.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/bore.log
stderr_logfile=/var/log/bore.log
priority=40
startsecs=3
startretries=999

[program:sftp-updater]
command=/usr/local/bin/update-sftp.sh
autostart=true
autorestart=true
priority=50

[program:keep-alive]
command=/usr/local/bin/keep-alive.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/keepalive.log
stderr_logfile=/var/log/keepalive.log
priority=60

[program:auto-cleanup]
command=/usr/local/bin/auto-cleanup.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/cleanup.log
stderr_logfile=/var/log/cleanup.log
priority=70

[program:screen-watchdog]
command=/usr/local/bin/screen-watchdog.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/screen-watchdog.log
stderr_logfile=/var/log/screen-watchdog.log
priority=80
