# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.log
priority=10
startsecs=3
startretries=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.log
priority=20

[program:ngrok]
command=ngrok tcp 22 --log stdout
autostart=true
autorestart=true
stdout_logfile=/var/log/tunnel.log
stderr_logfile=/var/log/tunnel.log
priority=30
startsecs=5
startretries=999

[program:bore-fallback]
command=/usr/local/bin/bore local 22 --to bore.pub
autostart=true
autorestart=true
stdout_logfile=/var/log/bore.log
stderr_logfile=/var/log/bore.log
priority=31
startsecs=5
startretries=3

[program:tunnel-updater]
command=/usr/local/bin/update-tunnel.sh
autostart=true
autorestart=true
priority=50
